# web/betterdesk.service
[Unit]
Description=BetterDesk Console - RustDesk Web Management
After=network.target

[Service]
Type=simple
User=root
WorkingDirectory=/opt/BetterDeskConsole
ExecStart=/usr/bin/python3 /opt/BetterDeskConsole/app.py
Restart=always
RestartSec=10

[Install]
WantedBy=multi-user.target
