# web/betterdesk.service
[Unit]
Description=BetterDesk Web Console v2.1.1
Documentation=https://github.com/UNITRONIX/Rustdesk-FreeConsole
After=network.target rustdesksignal.service

[Service]
Type=simple
User=root
WorkingDirectory=/opt/BetterDeskConsole
Environment="FLASK_HOST=0.0.0.0"
Environment="FLASK_PORT=5000"
Environment="FLASK_DEBUG=False"
Environment="RUSTDESK_PATH=/opt/rustdesk"
Environment="API_PORT=21120"
ExecStart=/opt/BetterDeskConsole/venv/bin/python /opt/BetterDeskConsole/app.py
Restart=always
RestartSec=10

[Install]
WantedBy=multi-user.target
